# etc/profile.d/firefox-touchsupport.sh
#
# Fix Firefox Touch Support
#
export MOZ_USE_XINPUT2=1
