# etc/systemd/system/bluetoothdisable.service
[Unit]
Description=Disable Bluetooth

[Service]
Type=oneshot
ExecStart=/usr/bin/rfkill block bluetooth

[Install]
WantedBy=multi-user.target
